# data/hebkb.conf
# Hebrew keyboard overlay. Before kb starts: bind kbroot (incl he/) over the
# read-only /usr/share/keyboard squashfs, then honor Keyboard.preferences
# (keyboard=he|en_US) by setting conf current. /usr/share/keyboard is always a
# mountpoint (squashfs) -> test for the he/ dir, never `mountpoint`.
start on starting kb
task
script
    if [ -f /var/local/kbroot/he/he-600x800.keymap.gz ]; then
        [ -d /usr/share/keyboard/he ] || mount --bind /var/local/kbroot /usr/share/keyboard
    fi
    KB=$(sed -n 's/^keyboard=//p' /var/local/java/prefs/Keyboard.preferences)
    if [ "$KB" = "he" ] && [ -d /usr/share/keyboard/he ]; then
        sed -i 's/"current": *"[^"]*"/"current": "he"/; s/"selected": *"[^"]*"/"selected": "he"/' /var/local/system/keyboard.conf
    fi
end script
