# data/hebkb.conf
# Hebrew keyboard overlay. Before kb starts: bind kbroot (incl he/) over the
# read-only /usr/share/keyboard squashfs, then honor Keyboard.preferences
# (keyboard=he|en_US) by setting conf current. /usr/share/keyboard is always a
# mountpoint (squashfs) -> test for the he/ dir, never `mountpoint`. The keymap
# guard is resolution-agnostic (any he-*.keymap.gz), and `selected` is kept as
# he:en_US so the on-screen globe key can switch Hebrew<->English.
start on starting kb
task
script
    if ls /var/local/kbroot/he/*.keymap.gz >/dev/null 2>&1; then
        [ -d /usr/share/keyboard/he ] || mount --bind /var/local/kbroot /usr/share/keyboard
    fi
    KB=$(sed -n 's/^keyboard=//p' /var/local/java/prefs/Keyboard.preferences)
    if [ "$KB" = "he" ] && [ -d /usr/share/keyboard/he ]; then
        sed -i 's/"current": *"[^"]*"/"current": "he"/; s/"selected": *"[^"]*"/"selected": "he:en_US"/' /var/local/system/keyboard.conf
    fi
end script
